# etc/profile.d/alias.sh
#!/usr/bin/bash
# System Management Daemon - Spaceport
# iDigitalFlame
#
# Default alias for nano.

alias nano="nano -Ll"
