# etc/profile.d/freetype2.sh
#!/usr/bin/bash
# System Management Daemon - Spaceport
# iDigitalFlame
#
# Fonts Profile Script.

export FREETYPE_PROPERTIES="truetype:interpreter-version=40"
